# 2.6/init-arm/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/busybox/bin/sh
::ctrlaltdel:/busybox/bin/umount -a -r
